# root/etc/services.d/homebridge/run
#!/usr/bin/with-contenv sh
sleep 7
s6-setuidgid abc homebridge -U /homebridge -P /homebridge/node_modules
