# root/etc/services.d/dbus/run
#!/usr/bin/with-contenv sh
dbus-daemon --system --nofork
