# rootfs/etc/systemd/system/homebridge.service.d/override.conf
# this docker container does not use systemd
# the homebridge apt package just checks this file to see if the user has been changed from the default

[Service]
User=root
