# root/etc/services.d/homebridge-config-ui-x/finish
#!/usr/bin/with-contenv sh

/usr/bin/container-credits-banner >> /homebridge/homebridge.log
