# root/etc/services.d/homebridge-config-ui-x/run
#!/usr/bin/with-contenv sh

s6-setuidgid root homebridge-config-ui-x \
  -U /homebridge \
  -P /homebridge/node_modules \
  -H /usr/local/lib/node_modules/homebridge \
  2>&1 | tee -a /homebridge/logs/homebridge.log
