# root/etc/services.d/avahi/run
#!/usr/bin/with-contenv sh
avahi-daemon
